# root/etc/services.d/snapcast/run
#!/usr/bin/with-contenv sh

# wait for services to come up
sleep 2

echo "Starting Snapcast Server..."

exec snapserver --config /etc/snapserver.conf